# voice-transcriber.service
[Unit]
Description=Voice Transcription Daemon
After=network.target sound.target
Wants=network.target sound.target

[Service]
Type=simple
User=%i
Group=%i
WorkingDirectory=/home/lukas/code/realtime-transcript-linux
ExecStart=/home/lukas/code/realtime-transcript-linux/venv/bin/python /home/lukas/code/realtime-transcript-linux/voice_daemon.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=DISPLAY=:0
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse

# Security settings
NoNewPrivileges=true
PrivateDevices=false
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/tmp /home/lukas/code/realtime-transcript-linux

# Resource limits
MemoryMax=1G
CPUQuota=50%

[Install]
WantedBy=default.target